# stt_rust/refiner-stt.service
[Unit]
Description=Refiner STT (Rust, native arm64)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=refiner
Group=refiner
WorkingDirectory=/opt/refiner/stt_rust
EnvironmentFile=-/etc/default/refiner-stt
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1
Environment=TOKIO_WORKER_THREADS=46
Environment=REFINER_STT_MODEL=/opt/refiner/models/ggml-base.en.bin
Environment=REFINER_STT_BIND=127.0.0.1:7079
Environment=REFINER_STT_LANG=en-GB
Environment=REFINER_STT_THREADS=2
Environment=REFINER_STT_WORKERS=23
Environment=REFINER_STT_MAX_AUDIO_BYTES=8000000
ExecStart=/opt/refiner/stt_rust/refiner-stt --model ${REFINER_STT_MODEL} --bind ${REFINER_STT_BIND} --lang ${REFINER_STT_LANG} --threads ${REFINER_STT_THREADS} --workers ${REFINER_STT_WORKERS} --max-audio-bytes ${REFINER_STT_MAX_AUDIO_BYTES}
CPUAffinity=0-45
TasksMax=8192
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
